# build/root/etc/supervisor/conf.d/supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
user=root
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n

[program:postfix]
command=/usr/sbin/postfix start-fg

[program:dovecot]
command=/usr/local/sbin/dovecot-start-wating-for-postfix.sh

[program:php8.4-fpm]
command=/usr/sbin/php-fpm8.4 --nodaemonize --fpm-config=/etc/php/8.4/fpm/pool.d/www.conf

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
